# com.plocate-macos.updatedb.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.plocate-macos.updatedb</string>

	<!--
	  require-visibility no: Skip per-user file permission checks on search
	  results. Safe for single-user Macs (the common case). On shared machines,
	  remove these two lines and create a _plocate group so each user only sees
	  paths they can access.
	-->
	<key>ProgramArguments</key>
	<array>
		<string>/opt/homebrew/sbin/updatedb</string>
		<string>--require-visibility</string>
		<string>no</string>
		<string>--output</string>
		<string>/opt/homebrew/var/lib/plocate/plocate.db</string>
	</array>

	<key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer>2</integer>
		<key>Minute</key>
		<integer>30</integer>
	</dict>

	<key>UserName</key>
	<string>root</string>

	<key>Nice</key>
	<integer>10</integer>

	<key>LowPriorityIO</key>
	<true/>

	<key>StandardOutPath</key>
	<string>/var/log/plocate-updatedb.log</string>

	<key>StandardErrorPath</key>
	<string>/var/log/plocate-updatedb.log</string>

	<key>EnvironmentVariables</key>
	<dict>
		<key>LANG</key>
		<string>en_US.UTF-8</string>
		<key>LC_ALL</key>
		<string>en_US.UTF-8</string>
	</dict>

	<key>RunAtLoad</key>
	<false/>
</dict>
</plist>
